# src/systemd/mapper-wait@.service.in
[Unit]
Description=Wait for %I
Wants=xyz.openbmc_project.ObjectMapper.service
After=xyz.openbmc_project.ObjectMapper.service

[Service]
Type=oneshot
RemainAfterExit=yes
Restart=on-failure
TimeoutStartSec=infinity
ExecStart=@BINDIR@/mapper wait %I
