# src/systemd/mapper-wait@.service.in
[Unit]
Description=Wait for %I

[Service]
Type=oneshot
RemainAfterExit=yes
Restart=on-failure
TimeoutStartSec=infinity
ExecStart=@BINDIR@/mapper wait %I
